# hab.service
[Unit]
Description=Samba
After=docker.service
Requires=docker.service

[Service]
EnvironmentFile=/etc/profile.env
TimeoutStartSec=0
ExecStartPre=-/usr/bin/docker kill openhab
ExecStartPre=-/usr/bin/docker rm openhab
ExecStart=/usr/bin/docker run --name openhab \
    --net=host \
    -p 8080 \
    --tty \
    -v ${MOUNT_ROOT}/openhab/openhab_addons:/openhab/addons \
    -v ${MOUNT_ROOT}/openhab/openhab_conf:/openhab/conf \
    -v ${MOUNT_ROOT}/openhab/openhab_userdata:/openhab/userdata \
    -d --restart=always \
    openhab/openhab:2.1.0-amd64

[Install]
WantedBy=multi-user.target
