# media-raid.mount
[Unit]
Description=RAID1

[Mount]
What=/dev/md0
Where=/media/raid
Type=ext4
Options=defaults,auto,noatime,exec

[Install]
WantedBy=multi-user.target
