# ansible/roles/logrotate/templates/logrotate.timer
# minetest.timer
[Unit]
Description=Runs logrotate for traefik

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h

[Install] 
WantedBy=timers.target
